# bootid-logger/log-bootid-on-boot.conf
# Copyright 2020 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Record the boot time and boot id."
author          "chromium-os-dev@chromium.org"

start on started syslog
task

exec /usr/sbin/bootid-logger
