# patchpanel/init/patchpanel.conf
# Copyright 2015 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Starts platform guest networking services"
author          "chromium-os-dev@chromium.org"

start on started nat
stop on stopped nat

respawn
respawn limit 3 10

pre-start script
{
  echo "Starting patchpaneld"
} 2>&1 | logger -t "${UPSTART_JOB}"
end script # pre-start

exec /usr/bin/patchpaneld

post-stop script
{
  echo "Stopping patchpaneld"
} 2>&1 | logger -t "${UPSTART_JOB}"
end script # post-stop
