# biod/study/init/syslog_fingerprint_study.conf
# Copyright 2020 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Setup the fingerprint study log file for syslog."
author        "chromium-os-dev@chromium.org"

# We must create the fingerprint study log file before syslog starts.
# This is because syslog cannot create this file itself, due to permissions.
start on starting syslog
task

env FPSTUDY_LOG_FILE=/var/log/fpstudy-audit.log

script
  touch "${FPSTUDY_LOG_FILE}"
  chown syslog "${FPSTUDY_LOG_FILE}"
end script
