# init/upstart/ext-pci-drivers-allowlist.conf
# Copyright 2020 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description "PCI drivers allowlist for external PCI (e.g. thunderbolt) devices"
author      "chromium-os-dev@chromium.org"

start on starting boot-services and started syslog

# This is a very small script with extremely small footprint. It finishes
# quickly (it's not a daemon), has a bare minimal memory foot print, but must
# be allowed to run and not be killed otherwise there shall be security
# vulnerabilities.
oom score never

task

exec /sbin/ext_pci_drivers_allowlist_init
