# init/upstart/tracefs-init.conf
# Copyright 2020 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Debugfs initialization script"
author          "chromium-os-dev@chromium.org"

# Initialize tracefs after the startup job (chromeos_startup script) finishes
# execution.
start on starting boot-services
# Use task job type to block other jobs that need tracefs.
task

# Disable OOM killer as we must never fail.
oom score never

script
  # Restore SELinux context for tracefs files.
  if [ -f /sys/fs/selinux/enforce ]; then
    restorecon -R /sys/kernel/debug/tracing
  fi
end script
