# arc/myfiles/arc-myfiles.conf
# Copyright 2016 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Mount /home/chronos/user/MyFiles for ARC."
author        "chromium-os-dev@chromium.org"

start on started arc-boot-continue
stop on stopped arc-boot-continue

env SOURCE=/home/chronos/user/MyFiles
env DEST=/run/arc/media/MyFiles
env UMASK=007
env ANDROID_MEDIA_UID=1023
env ANDROID_MEDIA_GID=1023

expect fork

pre-start script
  /usr/share/cros/init/logging-pre-start.sh

  exec logger -t "${UPSTART_JOB}" "Pre-start ${UPSTART_JOB}"
end script

exec /usr/share/cros/init/logging-exec.sh \
  mount-passthrough-jailed "${SOURCE}" "${DEST}" "${UMASK}" \
    "${ANDROID_MEDIA_UID}" "${ANDROID_MEDIA_GID}" "full"

post-start exec /usr/share/cros/init/logging-post-start.sh

post-stop exec logger -t "${UPSTART_JOB}" "Post-stop ${UPSTART_JOB}"
